# deployment/prompterly-backend.service
[Unit]
Description=Prompterly Backend API
After=network.target

[Service]
# Runs as the prompterly user so the app can read .env (chmod 600) and
# files under /home/prompterly/public_html (drwxr-x---, prompterly-owned).
# Previously User=www-data, which had no traversal rights into prompterly's
# home — pydantic-settings failed to open .env at startup.
User=prompterly
Group=prompterly
WorkingDirectory=/home/prompterly/app
Environment="PATH=/home/prompterly/app/.venv/bin"
EnvironmentFile=/home/prompterly/app/.env
ExecStart=/home/prompterly/app/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8001
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
