# deployment/prompterly-backend.service
[Unit]
Description=Prompterly Backend API
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/home/prompterly/public_html/api
Environment="PATH=/home/prompterly/public_html/api/.venv/bin"
EnvironmentFile=/home/prompterly/public_html/api/.env
ExecStart=/home/prompterly/public_html/api/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
